# resources/supervisord.conf
[supervisord]
nodaemon=true

[program:tomcat]
command=catalina.sh run
autorestart=false

[program:sshd]
command=/usr/sbin/sshd -D
